# usr/lib/nbld/nexenta_livecd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--
    Copyright 2005 Nexenta Systems, Inc.  All rights reserved.
    Use is subject to license terms.

    $Id: elatte_livecd.xml 10949 2005-11-20 03:48:09Z mac $
-->

<service_bundle type='profile' name='nexenta_livecd'
	 xmlns:xi='http://www.w3.org/2003/XInclude' >

  <!--
      svc.startd(1M) services
  -->
  <service name='system/dumpadm' version='1' type='service'>
    <instance name='default' enabled='false'/>
  </service>

  <service name='application/gdm' version='1' type='service'>
    <instance name='default' enabled='true'/>
  </service>

</service_bundle>
